# biome/themelio-node/scripts/ansible-debian-aws/files/update-images-and-restart.service
[Unit]
Description=Pull Latest Themelio-Node Image And Restart The Service
[Service]
User=admin
ExecStart=/usr/local/bin/update-images-and-restart.sh